# launchd/com.cc-reaper.guard-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
 <key>Label</key>
 <string>com.cc-reaper.guard-monitor</string>

 <key>ProgramArguments</key>
 <array>
 <string>/bin/bash</string>
 <string>__HOME__/.cc-reaper/cc-reaper-guard-monitor.sh</string>
 </array>

 <!-- Run every 30 seconds -->
 <key>StartInterval</key>
 <integer>30</integer>

 <key>RunAtLoad</key>
 <true/>

 <key>StandardOutPath</key>
 <string>__HOME__/.cc-reaper/logs/guard-launchd-stdout.log</string>

 <key>StandardErrorPath</key>
 <string>__HOME__/.cc-reaper/logs/guard-launchd-stderr.log</string>

 <key>EnvironmentVariables</key>
 <dict>
 <key>CC_MAX_SESSIONS</key>
 <string>3</string>
 <key>CC_IDLE_THRESHOLD</key>
 <string>1</string>
 <key>CC_MAX_RSS_MB</key>
 <string>3072</string>
 <key>CC_MAX_DESCENDANTS</key>
 <string>96</string>
 <key>CC_MAX_ZOMBIES</key>
 <string>16</string>
 <key>CC_MAX_GROWTH_MB_PER_MIN</key>
 <string>512</string>
 <key>CC_GUARD_MIN_GROWTH_MEM_MB</key>
 <string>1024</string>
 <key>CC_USE_FOOTPRINT</key>
 <string>1</string>
 <key>CC_GUARD_NOTIFY</key>
 <string>0</string>
 <key>CC_GUARD_KILL_BLOATED</key>
 <string>1</string>
 <key>CC_GUARD_KILL_DESCENDANTS</key>
 <string>1</string>
 <key>CC_GUARD_KILL_GROWTH</key>
 <string>1</string>
 </dict>
</dict>
</plist>
